# two-day-workshop/lab-cheat-sheets/lab4/step3/machine-2-hab-sup.service
[Unit]
Description=Habitat Supervisor

[Service]
ExecStart=/bin/hab sup run \
                           --listen-ctl 0.0.0.0:9632 \
                           --peer <machine-1-ip>

[Install]
WantedBy=default.target
